# settings/autostart/Compositor.desktop
[Desktop Entry]
Exec=picom --config ~/.config/picom/picom.conf
Name=Compositor
Type=Application
Version=1.0
