# settings/autostart/kdeconnect.desktop
[Desktop Entry]
Exec=kdeconnect-indicator
Name=kdeconnect
Type=Application
Version=1.0
