# settings/autostart/Night Light.desktop
[Desktop Entry]
Exec=redshift
Name=Night Light
Type=Application
Version=1.0
